# systemd/auto-brightness.service
[Unit]
Description=Auto-Brightness Daemon (Epilepsy-Safe)
After=graphical-session.target network.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/daemon
Restart=always
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
